# meme_bot.service
[Unit]
Description=Meme Coin Analysis Bot
After=network.target

[Service]
# Change this to your username (e.g., root or ubuntu)
User=root
# Change this to where you uploaded the code
WorkingDirectory=/root/meme-bot
# Path to python (run 'which python3' to confirm)
ExecStart=/usr/bin/python3 bot/main.py
Restart=always
RestartSec=10
# Environment variables (Optional if not in config.py)
# Environment=TELEGRAM_BOT_TOKEN=your_token

[Install]
WantedBy=multi-user.target
